# medishDevice.service
 
[Service]
Type=simple
ExecStart=/opt/MedishDevice/bin/MedishDevice

[Install]
WantedBy=multi-user.target
